# go/sound.conf
start on started httpd-user
stop on runlevel [06]

chdir /volume1/other/soundDir
exec ./sound