# go/sound.conf
start on syno.share.ready and syno.network.ready and started pgsql-adapter
stop on runlevel [06]

chdir /volume1/other/soundDir
exec ./sound
